# bitwarden-lite/supervisord/admin.ini
[program:admin]
autostart=%(ENV_BW_ENABLE_ADMIN)s
autorestart=true
command=/usr/bin/dotnet "Admin.dll"
directory=/app/Admin
environment=ASPNETCORE_URLS="http://+:5000"
priority=3
redirect_stderr=true
startsecs=15
stdout_logfile=/etc/bitwarden/logs/supervisord/admin.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
