# bitwarden-lite/supervisord/admin.ini
[program:admin]
autostart=true
autorestart=true
command=/app/Admin/Admin
directory=/app/Admin
environment=ASPNETCORE_URLS="http://+:5000"
priority=3
redirect_stderr=true
startsecs=15
stdout_logfile=/var/log/bitwarden/admin.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
