# bitwarden-lite/supervisord/api.ini
[program:api]
autostart=%(ENV_BW_ENABLE_API)s
autorestart=true
command=/usr/bin/dotnet "Api.dll"
directory=/app/Api
environment=ASPNETCORE_URLS="http://+:5001"
priority=2
redirect_stderr=true
startsecs=15
stdout_logfile=/etc/bitwarden/logs/supervisord/api.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
