# bitwarden-lite/supervisord/events.ini
[program:events]
autostart=%(ENV_BW_ENABLE_EVENTS)s
autorestart=true
command=/usr/bin/dotnet "Events.dll"
directory=/app/Events
environment=ASPNETCORE_URLS="http://+:5003"
priority=3
redirect_stderr=true
startsecs=15
stdout_logfile=/etc/bitwarden/logs/supervisord/bitwarden/events.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
