# bitwarden-lite/supervisord/identity.ini
[program:identity]
autostart=%(ENV_BW_ENABLE_IDENTITY)s
autorestart=true
command=/usr/bin/dotnet "Identity.dll"
directory=/app/Identity
environment=ASPNETCORE_URLS="http://+:5005"
priority=1
redirect_stderr=true
startsecs=15
stdout_logfile=/etc/bitwarden/logs/supervisord/identity.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
