# bitwarden-lite/supervisord/notifications.ini
[program:notifications]
autostart=%(ENV_BW_ENABLE_NOTIFICATIONS)s
autorestart=true
command=/usr/bin/dotnet "Notifications.dll"
directory=/app/Notifications
environment=ASPNETCORE_URLS="http://+:5006"
priority=3
redirect_stderr=true
startsecs=15
stdout_logfile=/etc/bitwarden/logs/supervisord/bitwarden/notifications.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
