# bitwarden-lite/supervisord/scim.ini
[program:scim]
autostart=%(ENV_BW_ENABLE_SCIM)s
autorestart=true
command=/usr/bin/dotnet "Scim.dll"
directory=/app/Scim
environment=ASPNETCORE_URLS="http://+:5002"
priority=4
redirect_stderr=true
startsecs=15
stdout_logfile=/etc/bitwarden/logs/supervisord/scim.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
