# bitwarden-lite/supervisord/sso.ini
[program:sso]
autostart=true
autorestart=true
command=/usr/bin/dotnet "Sso.dll"
directory=/app/Sso
environment=ASPNETCORE_URLS="http://+:5007"
priority=4
redirect_stderr=true
startsecs=15
stdout_logfile=/var/log/bitwarden/sso.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=5
